# modules/data/templates/udev/lowbat.rules.epp
# Suspend the system when battery level drops to <%= $profile::lowbat_suspend::low_percentage %>% or lower
SUBSYSTEM=="power_supply", ATTR{model_name}=="DELL JHXPY53", ATTR{status}=="Discharging", ATTR{capacity}=="[0-<%= $profile::lowbat_suspend::low_percentage %>]", RUN+="<%= $profile::lowbat_suspend::command %>"
